# scripts/com.knowledge-mcp.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.knowledge-mcp</string>

	<key>ProgramArguments</key>
	<array>
		<string>/Users/aq/knowledge-mcp/knowledge-mcp</string>
		<string>serve</string>
		<string>--mcp</string>
	</array>

	<key>EnvironmentVariables</key>
	<dict>
		<key>KNOWLEDGE_MCP_DATA_DIR</key>
		<string>/Users/aq/knowledge-mcp/data</string>
		<key>KNOWLEDGE_MCP_DEFAULT_KB</key>
		<string>notes</string>
		<key>KNOWLEDGE_MCP_LOG_FILE</key>
		<string>/Users/aq/knowledge-mcp/knowledge-mcp-launchd.log</string>
		<key>EMBED_API_ENDPOINT</key>
		<string>http://localhost:11434/v1/embeddings</string>
		<key>EMBED_MODEL</key>
		<string>bge-m3</string>
		<key>EMBED_DIM</key>
		<string>1024</string>
		<key>RERANK_API_ENDPOINT</key>
		<string>http://localhost:7997/rerank</string>
		<key>RERANK_CANDIDATE_LIMIT</key>
		<string>100</string>
	</dict>

	<!-- Automatically restart on crash -->
	<key>KeepAlive</key>
	<true/>

	<!-- Restart with a 10-second delay after crash -->
	<key>ThrottleInterval</key>
	<integer>10</integer>

	<!-- Log files -->
	<key>StandardOutPath</key>
	<string>/Users/aq/knowledge-mcp/knowledge-mcp-launchd.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/aq/knowledge-mcp/knowledge-mcp-launchd.log</string>

	<!-- Run at user login -->
	<key>RunAtLoad</key>
	<true/>

	<!-- Use project directory as working directory -->
	<key>WorkingDirectory</key>
	<string>/Users/aq/knowledge-mcp</string>
</dict>
</plist>
